# templates/cfupdate.service.j2
[Unit]
Description=Cloudflare DDNS service
After=network.target

[Service]
Type=oneshot
ExecStart={{ cfupdate_script_path }}/{{ cfupdate_script_filename }}

[Install]
WantedBy=multi-user.target
